# pi-web.service
[Unit]
Description=Pi Web App - AI Coding Agent (SDK in-process)
After=network.target

[Service]
Type=simple
WorkingDirectory=/home/manu/pi-web-app

# Graceful shutdown
ExecStop=/bin/bash -c 'kill -TERM $MAINPID; sleep 3; kill -9 $MAINPID 2>/dev/null; exit 0'
ExecStartPre=/bin/bash -c 'fuser -k 3210/tcp 2>/dev/null; sleep 2; exit 0'
ExecStart=/home/manu/.npm/_npx/fd45a72a545557e9/node_modules/.bin/tsx /home/manu/pi-web-app/src/server.ts

Restart=always
RestartSec=5

# Environment
Environment=NODE_ENV=production
Environment=PI_WEB_CWD=/home/manu
Environment=PI_WEB_PORT=3210
Environment=HOME=/home/manu
Environment=NVM_DIR=/home/manu/.nvm
Environment=PATH=/home/manu/.nvm/versions/node/v24.12.0/bin:/usr/local/bin:/usr/bin:/bin
Environment=NPM_CONFIG_PREFIX=/home/manu/.nvm/versions/node/v24.12.0

ReadWritePaths=/home/manu/.pi/agent /home/manu /tmp /home/manu/pi-web-app
LimitNOFILE=65536

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=pi-web

[Install]
WantedBy=default.target